# airvantage/packages/debian/postinst
#!/bin/bash
#*******************************************************************************
# Copyright (c) 2012 Sierra Wireless and others.
# All rights reserved. This program and the accompanying materials
# are made available under the terms of the Eclipse Public License v1.0
# and Eclipse Distribution License v1.0 which accompany this distribution.
#
# The Eclipse Public License is available at
#   http://www.eclipse.org/legal/epl-v10.html
# The Eclipse Distribution License is available at
#   http://www.eclipse.org/org/documents/edl-v10.php
#
# Contributors:
#     Guilhem Saurel for Sierra Wireless - initial script
#*******************************************************************************

useradd -b /opt -U mihini || echo "Mihini user already exists"
usermod -aG tty mihini
usermod -aG dialout mihini

cd ~mihini

for rw in persist update apps
do
    mkdir -p /var/opt/mihini/$rw
    ln -s /var/opt/mihini/$rw
done

chown -R mihini:mihini /{var/,}opt/mihini

cat << EOF > /etc/init.d/mihini
#!/bin/sh

### BEGIN INIT INFO
# Provides: mihini_agent
# Required-Start: \$remote_fs \$syslog
# Required-Stop: \$remote_fs \$syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Start mihini agent at boot time
# Description: Enable application framework or M2M services (like device management) provided by mihini agent.
### END INIT INFO

DESC="Eclipse Mihini M2M service"
NAME=Mihini
DAEMON=/opt/mihini/bin/appmon_daemon
AGENT=/opt/mihini/bin/agent
ROOTDIR=/opt/mihini
APPMON_PIDFILE=/run/mihini.appmon.pid
AGENT_PIDFILE=/run/mihini.agent.pid

. /lib/lsb/init-functions

test -f \$DAEMON || exit 0
test -f \$AGENT || exit 0

do_start()
{
    log_daemon_msg "Starting \$DESC" "Applications daemon"
    start-stop-daemon --start --quiet --pidfile \$APPMON_PIDFILE --exec \$DAEMON -- -w \$ROOTDIR -h mihini -u mihini -g mihini -n 5
    log_end_msg \$?
    test -e \$APPMON_PIDFILE || pidof appmon_daemon > \$APPMON_PIDFILE

    log_daemon_msg "Starting \$DESC" "Agent"
    start-stop-daemon --start --quiet --make-pidfile --pidfile \$AGENT_PIDFILE --chuid root:mihini --chdir \$ROOTDIR --background --exec \$AGENT
    log_end_msg \$?
}

do_stop()
{
    log_daemon_msg "Stopping \$DESC" "Agent"
    echo "destroy" | nc localhost 4242 > /dev/null 2>&1
    start-stop-daemon --stop --quiet --signal TERM --pidfile \$AGENT_PIDFILE
    log_end_msg \$?

    log_daemon_msg "Stopping \$DESC" "Applications daemon"
    start-stop-daemon --stop --quiet --oknodo --signal TERM --pidfile \$APPMON_PIDFILE 2>/dev/null
    log_end_msg \$?
    rm -f \$APPMON_PIDFILE \$AGENT_PIDFILE
}

case "\$1" in
    start)
        do_start
        ;;

    stop)
        do_stop
        ;;

    restart)
        do_stop
        do_start
        ;;

    status)
        status_of_proc \$DAEMON "\$NAME applications daemon"
        status_of_proc \$AGENT "\$NAME agent"
        ;;

    purge)
        rm -rf /var/opt/mihini/persist
        log_daemon_msg "Mihini's persist folder cleared"
        ;;

    *)
      echo "Usage: $0 {start|stop|restart|status|purge}"
      exit 1
      ;;
esac

exit 0
EOF

chmod +x /etc/init.d/mihini

echo "/opt/mihini/lib" > /etc/ld.so.conf.d/01-mihini.conf
ldconfig

/etc/init.d/mihini start
